# Egstart.service
[Service]
ExecStart=/bin/startm


[Install]
WantedBy=multi-user.target
